# files/etc/local.d/io-scheduler.start
#!/bin/bash
echo none > /sys/block/mmcblk0/queue/scheduler
