# files/etc/local.d/kms.start
/root/kms/arm/little-endian/static/vlmcsd-armv7el-uclibc-static &
